# deploy/systemd/deribit-sync-docker.service
[Unit]
Description=Deribit Options Data Sync Service (Docker)
After=network-online.target docker.service
Wants=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
User=sam
Group=sam
EnvironmentFile=/etc/downloader-sync.env

Environment=COMPOSE_PROJECT_NAME=deribit-sync
Environment=QUESTDB_HOST=localhost
Environment=QUESTDB_PORT=9009

# Sync BTC then ETH sequentially
ExecStart=/home/sam/.local/bin/dotenvx run -f /media/sam/1TB/.env -- /media/sam/1TB/monitoring-stack/scripts/cron-wrapper.sh deribit-sync-btc /usr/bin/docker compose --env-file /etc/downloader-sync.env --project-directory ${DERIBIT_REPO_ROOT} -f ${DERIBIT_REPO_ROOT}/docker-compose.yml run --rm deribit-sync
ExecStart=/home/sam/.local/bin/dotenvx run -f /media/sam/1TB/.env -- /media/sam/1TB/monitoring-stack/scripts/cron-wrapper.sh deribit-sync-eth /usr/bin/docker compose --env-file /etc/downloader-sync.env --project-directory ${DERIBIT_REPO_ROOT} -f ${DERIBIT_REPO_ROOT}/docker-compose.yml run --rm deribit-sync-eth

TimeoutStartSec=3600
Restart=no
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deribit-sync-docker

[Install]
WantedBy=multi-user.target
